# .docker/supervisor/hermes.ini
[supervisord]
logfile = /var/log/supervisor/supervisord.log
logfile_maxbytes = 10MB
logfile_backups=10
loglevel = info
nodaemon = false

[program:exchange-octane]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan octane:start --max-requests=1000 --workers=4 --task-workers=12 --port=9581
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=www-data
redirect_stderr=false

[program:exchange-horizon]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=www-data
redirect_stderr=false

[program:exchange-websockets]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan websockets:serve --port=6000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=root
redirect_stderr=false

[program:exchange-orders-update]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan orders:update
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=www-data
redirect_stderr=false

[program:exchange-markets-update]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan markets:update
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=www-data
redirect_stderr=false

[program:exchange-markets-update]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan blockchains:explore
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=www-data
redirect_stderr=false

[program:exchange-markets-update]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan accountancy:run
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=www-data
redirect_stderr=false

[program:exchange-kyc-authorization]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan kyc:authorization
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=0
user=www-data
redirect_stderr=false
